# etc/profile.d/desktopuser.sh

: "${CS:=dark}"
: "${CS_ADJ:=true}"

: "${X11_MONITOR_LAYOUT_ID:=${OS_HOSTNAME:?}-${OS_MODE:-desktop}}"

#
